# templates/bash-shell/change_desktop.service
[Unit]
Description=Change Desktop wallpaper

[Service]
User=/home/auri
ExecStart=%h/.opt/change_desktop

[Install]
WantedBy=multi-user.target
