# templates/bash-shell/change_desktop.service
[Unit]
Description=Change Desktop wallpaper

[Service]
ExecStart=%h/.opt/change_desktop

[Install]
WantedBy=default.target
